# deployment/rc.local
#!/bin/sh

/home/ubuntu/sorcerer/startup.sh
